# Dockerfile.d/etc/systemd/system/entrypoint.target
[Unit]
Description=target for entrypoint.service
Requires=entrypoint.service systemd-logind.service systemd-user-sessions.service